# application/configs/coredump-reporter.service
[Unit]
Description=deepin-log-viewer coredump report activities

[Service]

ExecStart=/usr/bin/deepin-log-viewer --reportcoredump
#CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYSLOG CAP_SYS_NICE CAP_SYS_PACCT CAP_MAC_OVERRIDE CAP_AUDIT_CONTROL CAP_AUDIT_READ CAP_AUDIT_WRITE
CapabilityBoundingSet=~CAP_NET_RAW
MemoryLimit=8G

[Install]
WantedBy=multi-user.target
